# systemd/gcal-sync.service
[Unit]
Description=Football Home — Google Calendar mirror + classifier (Slices 2 + 3)
Documentation=file:///srv/footballhome/docs/calendar-design.md
After=network-online.target
Wants=network-online.target
# Don't fire while the DB is being started/restarted by podman.
After=podman.service

[Service]
Type=oneshot
User=jbreslin
Group=jbreslin
WorkingDirectory=/srv/footballhome
# Two ExecStart lines run sequentially for Type=oneshot. Classifier
# only sees the events sync just wrote, so ordering matters. If sync
# fails, classifier is skipped (systemd stops on first non-zero exit).
#   1) gcal-sync.js     mirror Google → gcal_events
#   2) gcal-classify.js promote soccer-prefixed rows → fh_events
ExecStart=/usr/bin/node /srv/footballhome/scripts/gcal-sync.js
ExecStart=/usr/bin/node /srv/footballhome/scripts/gcal-classify.js
# Both scripts exit non-zero on Google-auth / DB failures; systemd will
# journal the failure but the .timer keeps firing every 5 min so a
# transient network hiccup self-heals on the next tick.
StandardOutput=journal
StandardError=journal
# Kill if it ever runs longer than 3 min — full sync of 1000+ events
# is ~50s + classifier <1s; incremental sync + classifier <2s total.
TimeoutStartSec=180s
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6

[Install]
WantedBy=multi-user.target
